# vector-input-data/gendl-sources-AGPL/docker/build-build/etc/rc.local
#!/bin/bash

/opt/gendl/run-gendl-services


